# scripts/ndi-web.service
[Unit]
Description=NDI Web Control Panel
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ndi-rpi-solution
Environment=PATH=/home/pi/ndi-env/bin:/usr/local/bin:/usr/bin
Environment=LD_LIBRARY_PATH=/opt/ndi/lib
ExecStart=/home/pi/ndi-env/bin/python3 /home/pi/ndi-rpi-solution/ndi_web.py --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
